# etc/supervisor.d/available/cjc.conf

[program:cjc]
command=%(ENV_HOME)s/.pyenv/shims/celery worker --app=lib.workers.cjc_tasks -l INFO -Q cjc -n 'cjc.%%h' --concurrency=1 -O fair -P solo
directory=/home/fishface/FishFace2
user=fishface
numprocs=1
stdout_logfile=/home/fishface/FishFace2/var/log/workers/cjc.log
stderr_logfile=/home/fishface/FishFace2/var/log/workers/cjc.err
autostart=true
autorestart=true
stopwaitsecs=600
killasgroup=true
priority=100