# etc/supervisor.d/available/results.conf

[program:results]
command=%(ENV_HOME)s/.pyenv/shims/celery worker --app=lib.workers.results_tasks -l INFO -Q results -n 'results.%%h' --concurrency=4 -O fair 
directory=/home/fishface/FishFace2
user=fishface
numprocs=1
stdout_logfile=/home/fishface/FishFace2/var/log/workers/results.log
stderr_logfile=/home/fishface/FishFace2/var/log/workers/results.err
autostart=true
autorestart=true
stopwaitsecs=600
killasgroup=true
priority=100